# packer/services/grafana-backup.timer
[Unit]
Description=Grafana Backup Daily

[Timer]
OnCalendar=*-*-* 00:30:00
Persistent=true

[Install]
WantedBy=timers.target
